# hack/lldb/usr/share/containers/systemd/alpine.container
[Container]
Image=docker.io/alpine:latest
